# files/ndp-proxy-go.in
#!/bin/sh

# PROVIDE: ndp_proxy_go
# REQUIRE: NETWORKING
# KEYWORD: shutdown

. /etc/rc.subr

name=ndp_proxy_go
rcvar=ndp_proxy_go_enable
desc="IPv6 Neighbor Discovery Proxy"

load_rc_config $name

: ${ndp_proxy_go_enable:=NO}
: ${ndp_proxy_go_upstream:=""}
: ${ndp_proxy_go_downstream:=""}
: ${ndp_proxy_go_flags:=""}

command="/usr/local/sbin/ndp-proxy-go"
pidfile="/var/run/${name}.pid"

start_precmd="ndp_proxy_go_precmd"

ndp_proxy_go_precmd()
{
    if [ -z "${ndp_proxy_go_upstream}" ]; then
        echo "Error: ndp_proxy_go_upstream not set"
        return 1
    fi
    if [ -z "${ndp_proxy_go_downstream}" ]; then
        echo "Error: ndp_proxy_go_downstream not set"
        return 1
    fi
}

command_args="${ndp_proxy_go_flags} ${ndp_proxy_go_upstream} ${ndp_proxy_go_downstream} > /var/log/${name}.log 2>&1 &"

run_rc_command "$1"
